# mkosi.profiles/droid-juicer/mkosi.extra/usr/lib/systemd/system/droid-juicer.service.d/skip-kcmdline.conf
[Unit]
ConditionKernelCommandLine=!droid-juicer.skip
