# mkosi.extra/usr/lib/systemd/system/bluetooth.service.d/10-persist-setup.conf
[Unit]
Wants=persist-setup.service
After=persist-setup.service
